# deploy/systemd/uav-yolo.service
[Unit]
Description=UAV YOLO vision and MJPEG stream
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/uav_project/uav_system-rk3588
ExecStart=%h/anaconda3/envs/yolo/bin/python -u yolo_app/main.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
