# deploy/systemd/uav-app.service
[Unit]
Description=UAV control application and Web UI
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/uav_project/uav_system-rk3588
ExecStart=%h/anaconda3/envs/app/bin/python -m app.main --send-commands false
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
